# systemd/exam-ready.service
[Unit]
Description=Bistec Exam Readiness Platform
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=exam-ready
Group=exam-ready
WorkingDirectory=/opt/exam-ready/app
EnvironmentFile=/etc/exam-ready.env
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=exam-ready

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/exam-ready

[Install]
WantedBy=multi-user.target
